# nut/rootfs/etc/services.d/powerman/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Network UPS Tools
# Monitor PowerMan daemon for PDU support
# ==============================================================================

# Only run in netserver mode
if ! bashio::config.equals 'mode' 'netserver' ;then
    bashio::log.debug "Not in netserver mode, PowerMan monitor not needed"
    exec sleep 864000
fi

# Check if PowerMan should be running
has_powerman=false

# Check explicit enable flag
if bashio::config.has_value 'powerman_enabled' && bashio::config.true 'powerman_enabled' ;then
    has_powerman=true
fi

# Check for powerman-pdu devices
for device in $(bashio::config "devices|keys"); do
    driver=$(bashio::config "devices[${device}].driver")
    if [[ "${driver}" == "powerman-pdu" ]]; then
        has_powerman=true
        break
    fi
done

if [[ "${has_powerman}" == "false" ]]; then
    bashio::log.debug "No PowerMan PDU devices configured, PowerMan monitor not needed"
    exec sleep 864000
fi

# PowerMan should already be started by 00-powerman.sh init script
# This service just monitors it and restarts if needed

bashio::log.info "Monitoring PowerMan daemon..."

# Check if PowerMan is already running
if pgrep powermand > /dev/null; then
    bashio::log.info "PowerMan daemon already running, monitoring..."
    # Just sleep and let s6 restart us if powerman dies
    exec sleep 864000
fi

# If we get here, PowerMan isn't running but should be
bashio::log.warning "PowerMan daemon not found, attempting to start..."

# Check if configuration exists
if [ ! -f /etc/powerman/powerman.conf ]; then
    bashio::log.error "PowerMan configuration file not found!"
    exec sleep 864000
fi

# Try to start PowerMan
exec /usr/sbin/powermand -f -c /etc/powerman/powerman.conf
