# nut/rootfs/etc/services.d/powerman/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Network UPS Tools
# Run PowerMan daemon for PDU support
# ==============================================================================
if ! bashio::config.equals 'mode' 'netserver' ;then
    exec sleep 864000
fi

# Only start if powerman-pdu driver is configured
if bashio::config.contains 'devices[].driver' 'powerman-pdu' ;then
    bashio::log.info "Starting PowerMan daemon for PDU support..."
    exec powermand -f -c /etc/powerman/powerman.conf
else
    exec sleep 864000
fi
