# nut/rootfs/etc/services.d/powerman/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Network UPS Tools
# PowerMan service disabled - PowerMan is started in init script
# ==============================================================================

# This service is disabled because PowerMan is started in the init script
# /etc/cont-init.d/00-powerman.sh to ensure proper startup order

bashio::log.info "PowerMan service disabled (started in init script)"
exec sleep 864000
