# nut/rootfs/etc/services.d/powerman/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Network UPS Tools
# Run PowerMan daemon for PDU support
# ==============================================================================
if ! bashio::config.equals 'mode' 'netserver' ;then
    exec sleep 864000
fi

# Only start if PowerMan is enabled
if bashio::config.has_value 'powerman_enabled' && bashio::config.true 'powerman_enabled' ;then
    bashio::log.info "Starting PowerMan daemon for PDU support..."
    exec powermand -f -c /etc/powerman/powerman.conf
else
    exec sleep 864000
fi
